# .deploy/config/supervisor.conf
[supervisord]
logfile=/dev/null
logfile_maxbytes=0
logfile_backups=0
loglevel=info
nodaemon=true

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/srv/app/app
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=1

[program:celery-worker]
command=celery -A app worker --loglevel=%(ENV_LOG_LEVEL)s
directory=/srv/app/app
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=1

[program:celery-beat]
command=celery -A app beat --loglevel=%(ENV_LOG_LEVEL)s
directory=/srv/app/app
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=1

[program:flower]
command=celery -A app flower --port=5555 --address=0.0.0.0 --url_prefix=flower
directory=/srv/app/app
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PYTHONUNBUFFERED=1

[group:celery]
programs=celery-worker,celery-beat,flower
priority=999